# deps/supervisord/supervisord.conf
[supervisord]
loglevel=debug

[program:hitch]
command=/app/hitch/start.sh
user=hitch

[program:maracuya]
command=/app/maracuya/bin/server.js
user=maracuya